# ai_diagbe.service
[Unit]
Description=AI Diagnosis Backend Service
After=network.target

[Service]
User=root
WorkingDirectory=/opt/ai_diagbe
ExecStart=/opt/ai_diagbe/.venv/bin/python3 -m app.main
Restart=always
Environment=PYTHONPATH=/opt/ai_diagbe

[Install]
WantedBy=multi-user.target
